# getty@tty1.service.d/autologin.conf
[Service]
# network doesn't load and X11 errors sometimes
#Type=simple
#Environment=XDG_SESSION_TYPE=x11
ExecStart=
ExecStart=-/sbin/agetty --noreset --noclear --autologin me - ${TERM}
